# ovl/systemd/default/etc/systemd/system/multi-user.target
[Unit]
Description=Multi-user

[Install]
WantedBy=default.target
