# ovl/attic/systemd/debug/etc/systemd/system/default.target

[Unit]
Description=Default
After=network.target
Requires=network.target
