# ovl/kubernetes/default/etc/systemd/system/kube-node.service
[Unit]
Description=Kubernetes Node
After=kube-prep.service

[Service]
Type=forking
ExecStart=/etc/systemd/bin/kube-node
TimeoutSec=10
StandardOutput=tty
StandardError=inherit

[Install]
WantedBy=multi-user.target
